# home/bazz/.config/systemd/user/cleanup-files.timer
[Unit]
Description=Weekly clean-up of the Downloads folder

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3h

[Install]
WantedBy=timers.target
