# home/bazz/.config/systemd/user/rsync-home.service
[Unit]
Description=Rsync my home dir weekly

[Service]
Type=oneshot
ExecStart=/home/bazz/.local/bin/backup_home
