# home/bazz/.config/systemd/user/rsync-home.timer
[Unit]
Description=rsync-home

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
