# etc/systemd/system/etc-backup.timer
[Unit]
Description=Run weekly /etc/ directory backup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
