# etc/systemd/system/host-backup.timer
[Unit]
Description=Run weekly local backup

[Timer]
OnCalendar=weekly
Persistent=false
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
