# etc/systemd/system/host-backup.service
[Unit]
Description=Weekly backup service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup_host

[Install]
WantedBy=multi-user.target
